# dotfiles/profile.d/bats
#!/bin/bash -e

if [ -d "$HOME/.bats" ]; then
  export PATH="$HOME/.bats/bin:$PATH"
fi
