# dotfiles/profile.d/development
#!/bin/sh

code() {
  VSCODE_CWD="$PWD" open -nW -a "Visual Studio Code" --args "$@"
}

vscode() {
  code "$@"
}

alias s="start_dev_server"

start_dev_server() {
  if [ -f ./Procfile.dev ]; then
    echo "Procfile.dev detected"
    if [ -f ./config/boot.rb ]; then
      echo "Removing logs..."
      bundle exec rails log:clear
    fi
    echo "Running: foreman start -f Procfile.dev"
    foreman start -f Procfile.dev
  elif [ -f ./Procfile ]; then
    echo "Procfile detected"
    if [ -f ./config/boot.rb ]; then
      echo "Removing logs..."
      bundle exec rails log:clear
    fi
    echo "Running: foreman start"
    foreman start
  elif [ -f ./package.json ]; then
    echo "NPM package detected"
    echo "Running: npm start"
    npm start
  else
    echo "\033[1;31m** No app server found. **"
    echo "\033[1;31m** Looking for a Procfile(.dev) or package.json file. **"
    return 1
  fi
}
