# dotfiles/profile.d/development
#!/bin/bash -e

code() {
  VSCODE_CWD="$PWD" open -nW -a "Visual Studio Code" --args "$@"
}

vscode() {
  code "$@"
}

alias s='start'

start() {
  local red="\033[1;31m"
  local purple="\033[1;35m"
  local green="\033[1;32m"
  local reset="\033[0;0m"
  local messagea
  local commanda

  if [ -f ./Procfile.dev ]; then
    messagea="Procfile.dev detected"
    commanda="foreman start -f Procfile.dev"
  elif [ -f ./Procfile ]; then
    messagea="Procfile detected"
    commanda="foreman start"
  elif [ -f ./package.json ]; then
    messagea="NPM app detected"
    commanda="npm start"
  else
    messagea="No app found. Looking for Procfile or package.json."
    echo -e "$red$messagea$reset"
    return 1
  fi

  echo -e "$purple$messagea$reset"
  if [ -n "$commanda" ]; then
    if [ -f ./config/boot.rb ]; then
      echo "Removing logs..."
      bundle exec rails log:clear
    fi

    echo -e "Running $green$commanda$reset"
    $commanda
  fi
}

clean_xcode() {
  rm -rf ~/Library/Developer/Xcode/DerivedData/*
}

pause_and_prompt_to_continue() {
  read -r -s -p "Press Enter key to continue..."
  echo ''
}
