# dotfiles/profile.d/development
#!/bin/sh

# Run `code` to open Visual Studio Code from the Terminal.
# https://code.visualstudio.com/docs/setup/mac#_launching-from-the-command-line
if [ -d "/Applications/Visual Studio Code.app/Contents/Resources/app/bin" ]; then
  export PATH="/Applications/Visual Studio Code.app/Contents/Resources/app/bin:$PATH"
fi

export EDITOR="code --wait"

alias s="start-development"
