# dotfiles/profile.d/git
#!/bin/sh

alias g='git'
alias gb='git branch'
alias gco='git checkout'
alias gfo='git fetch origin'
alias gfu='git fetch upstream'
alias gs='git status'
