# dotfiles/profile.d/google
#!/bin/sh

source_google_cloud_sdk_file() {
  if [ -f "/usr/local/Caskroom/google-cloud-sdk/latest/google-cloud-sdk/$1" ]; then
    # shellcheck disable=SC1090
    . "/usr/local/Caskroom/google-cloud-sdk/latest/google-cloud-sdk/$1"
  fi
}

if [ -n "$ZSH_VERSION" ]; then
  # Updates PATH for the Google Cloud SDK.
  source_google_cloud_sdk_file "path.zsh.inc"

  # Enables shell command completion for gcloud.
  source_google_cloud_sdk_file "completion.zsh.inc"
elif [ -n "$BASH" ]; then
  source_google_cloud_sdk_file "path.bash.inc"
  source_google_cloud_sdk_file "completion.bash.inc"
fi

unset -f source_google_cloud_sdk_file
