# dotfiles/profile.d/google
#!/bin/sh

caskroom_prefix="$(brew --caskroom)"

source_file() {
  if [ -f "$caskroom_prefix/google-cloud-sdk/latest/google-cloud-sdk/$1" ]; then
    # shellcheck disable=SC1090
    . "$caskroom_prefix/google-cloud-sdk/latest/google-cloud-sdk/$1"
  fi
}

[ -n "$ZSH_VERSION" ] && source_file "path.zsh.inc"
[ -n "$BASH" ] && source_file "path.bash.inc"

unset -f source_file
unset caskroom_prefix
