# dotfiles/profile.d/heroku
#!/bin/bash -e

alias hrc='heroku run rails console'
alias hrr='heroku run rake'
alias hrdm='heroku run rake db:migrate'
alias hp='heroku ps'
