# dotfiles/profile.d/node
#!/bin/sh

if [ -d "$HOME/.nodenv" ]; then
  eval "$(nodenv init -)"
fi
