# dotfiles/profile.d/system
#!/bin/sh

if [ -d "$HOME/.bin" ]; then
  export PATH="$HOME/.bin:$PATH"
fi

alias opendir='open -R .'

alias ..='cd ..'
alias home='cd ~'
alias documents='cd ~/Documents'
alias p='cd ~/Projects'
alias projects='cd ~/Projects'
alias laptop='cd ~/Projects/laptop'
alias dotfiles='cd ~/Projects/laptop/dotfiles'
alias jj='cd ~/Projects/javierjulio'

alias r='reload'
alias c='clear'
alias q='exit'

# shellcheck disable=SC1091
reload() {
  [ -n "$ZSH_VERSION" ] && . "$HOME/.zshrc"
  [ -n "$BASH" ] && . "$HOME/.profile"
}

executable() {
  chmod +x "$@"
}

search() {
  find . -type f -print0 | xargs -0 grep -l "$1"
}

# Empty the Trash on all mounted volumes and the main HDD
# Also, clear Apple’s System Logs to improve shell startup speed
emptytrash() {
  sudo rm -rfv /Volumes/*/.Trashes
  sudo rm -rfv ~/.Trash
  sudo rm -rfv /private/var/log/asl/*.asl
}
